# overrides/home/joonas/.config/systemd/user/backup-rsync.service
[Unit]
Description=rsync backups
OnFailure=failure-notification@%n

[Service]
Type=oneshot
ExecStart=sudo jsys backup --refresh-excludes

