# overrides/home/joonas/.config/systemd/user/backup-rsync.timer
[Unit]
Description=Run backup-rsync

[Timer]
OnBootSec=3h
OnUnitActiveSec=3h

[Install]
WantedBy=timers.target

