# overrides/home/joonas/.config/systemd/user/lowdiskspace-checker.service
[Unit]
Description=Low disk space checker
OnFailure=failure-notification@%n

[Service]
Type=oneshot
ExecStart=/usr/bin/jsys lowdiskspace-checker
