# overrides/home/joonas/.config/systemd/user/notifymerge.socket
[Socket]
ListenStream=/run/user/%U/org.freedesktop.Notifications.Notify_merge.sock

[Install]
WantedBy=sockets.target
