# overrides/home/joonas/.config/systemd/user/notifymerge.service
[Unit]
Description=Desktop notification merger

[Service]
ExecStart=/usr/bin/hautomo-client notification-merge

[Install]
WantedBy=multi-user.target
