# overrides/home/joonas/.config/systemd/user/varasto-fuse.service
[Unit]
Description=Varasto FUSE server

[Install]
WantedBy=default.target

[Service]
ExecStart=/usr/bin/sto fuse serve
Restart=always
RestartSec=10s
