# overrides/etc/systemd/system/mnt-varasto-bender.mount
[Unit]
ConditionHost=htpc

[Mount]
What=LABEL=bender
Where=/mnt/varasto/bender
Options=defaults,nofail,noatime

[Install]
WantedBy=multi-user.target
