# overrides/etc/systemd/system/mnt-varasto-fry.mount
[Unit]
ConditionHost=htpc

[Mount]
What=LABEL=fry
Where=/mnt/varasto/fry
Options=defaults,nofail,noatime

[Install]
WantedBy=multi-user.target
