# overrides/etc/systemd/system/mnt-varasto-amy.mount
[Unit]
ConditionHost=htpc

[Mount]
What=LABEL=amy
Where=/mnt/varasto/amy
Options=defaults,nofail,noatime

[Install]
WantedBy=multi-user.target
