# overrides/etc/systemd/system/mnt-varasto-calculon.mount
[Unit]
ConditionHost=htpc

[Mount]
What=LABEL=Calculon
Where=/mnt/varasto/calculon
Options=defaults,nofail,noatime

[Install]
WantedBy=multi-user.target
