# overrides/etc/systemd/system/mnt-varasto-leela.mount
[Unit]
ConditionHost=htpc

[Mount]
What=LABEL=Leela
Where=/mnt/varasto/leela
Options=defaults,nofail,noatime

[Install]
WantedBy=multi-user.target
