# overrides/usr/lib/systemd/system/tailscaled.service.d/10-state-redirect.conf
# https://weeknotes.barrucadu.co.uk/notes/186.html
# https://github.com/tailscale/tailscale/issues/3039

[Service]
BindPaths=/sysroot/apps/tailscale:/var/lib/tailscale
